# main/hri/hri-autostart.service
[Unit]
Description=HRI Human-Robot Interface Service (Direct Output)
After=multi-user.target
Wants=multi-user.target

[Service]
Type=simple
User=jetson
Group=jetson
WorkingDirectory=/home/jetson/workspace/vehicle_wbt/main/hri
Environment="DISPLAY=:0"
Environment="QT_QPA_PLATFORM=xcb"
Environment="QT_XCB_GL_INTEGRATION=none"
Environment="QT_OPENGL=software"
Environment="QT_QUICK_BACKEND=software"
Environment="QT_SCALE_FACTOR=1"
Environment="QT_AUTO_SCREEN_SCALE_FACTOR=0"
Environment="QT_SCREEN_SCALE_FACTORS=1"
Environment="QT_FONT_DPI=96"
Environment="QT_LOGGING_RULES=*=false"
ExecStart=/home/jetson/workspace/vehicle_wbt/main/hri/hri_startup.sh
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal
KillMode=mixed
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target